# containers/debian-base/systemd/system-update.timer
[Unit]
Description=Run system update weekly

[Timer]
OnCalendar=weekly
Persistent=true
Unit=system-update.service

[Install]
WantedBy=timers.target
